# Boot and Networking files/PV Powertech/ems_v3_boot_network_gateway_20260717_150954/etc/systemd/system/ems-v3-soc-controller.service.d/10-after-gateway.conf
[Unit]
After=ems-v3-gateway.service
Wants=ems-v3-gateway.service
