# ems_gateway_test_imx93/ems_networking_bootup_scripts_v2/systemd/nb-ems-gateway.service
[Unit]
Description=NorthBound EMS Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kinetics-grid-ems/northbound_ems_gateway
Environment=PYTHONPATH=/root/kinetics-grid-ems/northbound_ems_gateway/src
ExecStart=/usr/bin/python3 -m nb_ems_gateway.main --config configs/development.json --mock
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target