# northbound_ems_gateway/deployment/systemd/nb-ems-soc-controller.service
[Unit]
Description=Northbound EMS SOC-only field controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/kinetics-grid-ems/northbound_ems_gateway
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/kinetics-grid-ems/northbound_ems_gateway/tools/soc_only_controller.py --live --interval 5
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
